# conf/supervisor.ini
[program:projectname]
command = /root/.virtualenvs/projectname_env/bin/gunicorn core.wsgi:application -c /srv/projectname-project/conf/production.py
directory = /srv/projectname-project/gordon/
autostart = true
autorestart = true
stderr_logfile = /var/log/projectname-project/supervisor/errors.log
stdout_logfile = /var/log/projectname-project/supervisor/access.log
